# pinchy.service
[Unit]
Description=Pinchy – lightweight AI agent platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pinchy
Group=pinchy
WorkingDirectory=/opt/pinchy
ExecStart=/opt/pinchy/pinchy
Restart=always
RestartSec=5
Environment=PINCHY_HOME=/opt/pinchy/.pinchy
Environment=RUST_LOG=info
Environment=TOKIO_WORKER_THREADS=2

# Resource limits (Pi-friendly)
MemoryMax=512M
CPUQuota=75%
IOWeight=50
TasksMax=64

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pinchy/.pinchy
PrivateTmp=true

[Install]
WantedBy=multi-user.target
